# rootfs/etc/s6/postgres/run
#!/bin/bash

if [ -f ./setup ]
then
  source ./setup
fi

su-exec postgres postgres
